# scripts/keepassa-supervisord.ini
[program:keepassa]
stdout_logfile=/opt/keepassa/logs/stdout.log
command=/usr/local/jdk-1.8.0/bin/java -jar /opt/keepassa/KeepassaServer.jar
directory=/opt/keepassa
autostart=true
autorestart=true
redirect_stderr=true